# etc/zfs-check.cron.daily
#!/bin/sh

which zfs-check > /dev/null && zfs-check -sl
